# scripts/systemd/agent-hub-backend.service
[Unit]
Description=agent-hub backend (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/agent-hub/backend
Environment="PATH=%h/.local/bin:%h/agent-hub/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
EnvironmentFile=-%h/.env.local
ExecStart=%h/agent-hub/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8003
Restart=always
RestartSec=10

# Zombie prevention - kill entire cgroup on stop
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-hub-backend

[Install]
WantedBy=default.target
